# deploy/lightsail/templates/openclaw-budget-guard.timer
[Unit]
Description=Run OpenAI budget guard every 15 minutes

[Timer]
OnBootSec=3m
OnUnitActiveSec=15m
AccuracySec=1m
Persistent=true

[Install]
WantedBy=timers.target
